# backup/2026-02-25_154932/tower-dashboard.service.d/99-final.conf
[Service]
ExecStartPost=/bin/bash -lc /usr/local/bin/patch_reports_final.sh || true
